# crontabs/root
*/1 *   * * *   /usr/sbin/chk-linode.sh
